# auto-trader.service
[Unit]
Description=Auto Crypto Trading Journal v2 (Reaction Network)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=trader
Group=trader
WorkingDirectory=/opt/Auto-Crypto-Tradingjournal
ExecStart=/opt/Auto-Crypto-Tradingjournal/venv/bin/python main.py --paper
ExecStopPost=/opt/Auto-Crypto-Tradingjournal/scripts/backup_db.sh
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/opt/Auto-Crypto-Tradingjournal/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/Auto-Crypto-Tradingjournal /opt/Auto-Crypto-Tradingjournal/logs /opt/Auto-Crypto-Tradingjournal/data /opt/Auto-Crypto-Tradingjournal/backups

# Logging
StandardOutput=append:/opt/Auto-Crypto-Tradingjournal/logs/daemon-stdout.log
StandardError=append:/opt/Auto-Crypto-Tradingjournal/logs/daemon-stderr.log

[Install]
WantedBy=multi-user.target