# trading-journal.service
[Unit]
Description=Crypto Trading Journal
After=network.target

[Service]
Type=simple
# Replace <your-user> with your Linux username and update the paths accordingly
User=<your-user>
WorkingDirectory=/home/<your-user>/trading-journal
EnvironmentFile=/home/<your-user>/trading-journal/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
