# scripts/serviced/auto-trader.service
[Unit]
Description=Auto-Trader
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jmh
Group=jmh
WorkingDirectory=/opt/Auto-Trader
ExecStart=/opt/Auto-Trader/venv/bin/python main.py --paper
ExecStopPost=/opt/Auto-Trader/scripts/backup_db.sh
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/opt/Auto-Trader/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/Auto-Trader /opt/Auto-Trader/logs /opt/Auto-Trader/data /opt/Auto-Trader/backups

# Logging
StandardOutput=append:/opt/Auto-Trader/logs/daemon-stdout.log
StandardError=append:/opt/Auto-Trader/logs/daemon-stderr.log

[Install]
WantedBy=multi-user.target