# scripts/serviced/auto-trader-backtest.service
[Unit]
Description=Auto Trader Backtest
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jmh
Group=jmh
WorkingDirectory=/opt/Auto-Trader
ExecStart=/opt/Auto-Trader/venv/bin/python -m core.replay_driver --start 2026-01-01 --strategy paper_v3_learning_test
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/opt/Auto-Trader/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/Auto-Trader /opt/Auto-Trader/logs /opt/Auto-Trader/data /opt/Auto-Trader/backups /opt/Auto-Trader/temp

# Memory protection
MemoryMax=3G

# Logging
StandardOutput=append:/opt/Auto-Trader/logs/backtest-stdout.log
StandardError=append:/opt/Auto-Trader/logs/backtest-stderr.log

[Install]
WantedBy=multi-user.target
